# supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true

[program:dovecot]
command=/app/start.sh
autostart=true
autorestart=true
